# src/etc/crontabs/offlineimap
0 * * * * offlineimap -c /etc/offlineimap/offlineimap.conf -u basic
